# packages/marcopolo/marcopolo-0.1.4.tar.gz/marcopolo-0.1.4/daemon/python3/polod.service
[Unit]
Description=A simple service discovery protocol advertiser
Requires=network.target network-online.target
After=network.target network-online.target 
After=sys-subsystem-net-devices-net0.device

[Service]
ExecStart=/usr/bin/polod &>/dev/null
ExecReload=/usr/bin/poloreload.py &> /dev/null

Type=simple
User=root
Group=root
PIDFile=/var/run/polod.pid

[Install]
WantedBy=multi-user.target
