# packages/zymbit/zymbit-0.4.22.tar.gz/zymbit-0.4.22/files/arduino/etc/init.d/zymbit
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=99

export HOME=/

BIN="/usr/bin/zymbit"
DEFAULT=/etc/default/tunnel
OPTIONS="connector"

start() {
  env >> /tmp/env
	[ -f $DEFAULT ] && . $DEFAULT
	$BIN $OPTIONS > /tmp/zymbit.log &
}

stop() {
	kill -9 $(ps | grep -v grep | grep $BIN | awk '{print $1}')
}

reload() {
  restart
}

restart() {
  stop
  start
}

