# packages/zymbit/zymbit-0.4.1.tar.gz/zymbit-0.4.1/files/raspberrypi-arch/etc/systemd/system/zymbit.service
[Install]
WantedBy=multi-user.target

[Unit]
Description=Zymbit Cloud Connector
After=network.target

[Service]
ExecStart=/opt/zymbit/latest/bin/zymbit --logfile /var/log/zymbit connector
RestartSec=1
