# packages/marcopolo/marcopolo-0.1.6.tar.gz/marcopolo-0.1.6/daemon/python3/marcod.service
[Unit]
Description=A simple service discovery protocol resolver
Requires=network.target network-online.target
After=network.target network-online.target 
After=sys-subsystem-net-devices-net0.device

[Service]
ExecStart=/usr/bin/marcod &>/dev/null

Type=simple
User=root
Group=root
PIDFile=/var/run/marcod.pid

[Install]
WantedBy=multi-user.target