# packages/marcobootstrap.backend/marcobootstrap.backend-0.0.3.tar.gz/marcobootstrap.backend-0.0.3/daemon/systemd/marcomanager.service
[Unit]
Description=MarcoManager service
After=marco.service
After=polo.service

[Service]

ExecStart=/usr/bin/marcobootstrapbackendd &>/dev/null

Type=simple
User=root
Group=root
PIDFile=/var/run/marcobootstrapbackendd.pid

[Install]
WantedBy=multi-user.target
